# backend/config/systemd/siem-analyzer.service
[Unit]
Description=SIEM Event Analyzer Service
After=network.target postgresql.service redis.service elasticsearch.service
Requires=postgresql.service redis.service elasticsearch.service

[Service]
Type=simple
User=siem
Group=siem
WorkingDirectory=/opt/siem
ExecStart=/usr/local/bin/python3 -m analyzer.main
Restart=always
RestartSec=3
Environment=PYTHONPATH=/opt/siem
Environment=CONFIG_PATH=/etc/siem/analyzer.yml
Environment=MODEL_PATH=/opt/siem/models

[Install]
WantedBy=multi-user.target
