# backend/config/systemd/siem-dashboard.service
[Unit]
Description=SIEM Dashboard Service
After=network.target siem-api.service
Requires=siem-api.service

[Service]
Type=simple
User=siem
Group=siem
WorkingDirectory=/opt/siem/dashboard
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=3
Environment=NODE_ENV=production
Environment=PORT=8443
Environment=API_URL=https://localhost:8000
Environment=SSL_KEY=/etc/siem/ssl.key
Environment=SSL_CERT=/etc/siem/ssl.crt

[Install]
WantedBy=multi-user.target
